# scripts/kiosk/kiosk-wayland.service
# Systemd service for Chrome Kiosk with Wayland (cage compositor)
# Install to: /etc/systemd/system/kiosk.service
#
# This is an alternative to the X11-based setup that uses the
# lightweight "cage" Wayland compositor for a true kiosk experience.
#
# Install cage: sudo pacman -S cage
#
# Enable with:
#   sudo systemctl enable kiosk.service
#   sudo systemctl start kiosk.service
#
# IMPORTANT: Replace YOUR_USERNAME and YOUR_UID below!

[Unit]
Description=Chrome Kiosk (Wayland/Cage)
After=systemd-user-sessions.service
Wants=systemd-user-sessions.service
# If your home automation server runs as a systemd service, add:
# After=home-automation.service

[Service]
Type=simple
User=YOUR_USERNAME
# Get your UID with: id -u YOUR_USERNAME
Environment=XDG_RUNTIME_DIR=/run/user/YOUR_UID
Environment=KIOSK_URL=http://localhost:3000/dashboard#home-kiosk

# Cage options:
#   -s: Allow VT switching (Ctrl+Alt+F1-F6)
#   -d: Enable debug output (remove for production)
#
# Chromium options for kiosk mode
ExecStart=/usr/bin/cage -s -- /usr/bin/chromium \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --disable-session-crashed-bubble \
    --no-first-run \
    --disable-translate \
    --disable-features=TranslateUI \
    --check-for-update-interval=31536000 \
    --ozone-platform=wayland \
    ${KIOSK_URL}

# Auto-restart on crash (handles SIGILL, SIGSEGV, etc.)
Restart=always
RestartSec=3

# Allow some time for the display to be available
ExecStartPre=/bin/sleep 2

[Install]
WantedBy=multi-user.target
